# com.devforge.tech-reading-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.devforge.tech-reading-push</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/keita/Develop/DevForge/tech-daily-reading-web/auto-push.sh</string>
  </array>
  <!-- Cowork の読み物タスク(13:30＋jitter)の後。14:00 に push -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>14</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/keita/Develop/DevForge/tech-daily-reading-web/.auto-push.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/keita/Develop/DevForge/tech-daily-reading-web/.auto-push.log</string>
</dict>
</plist>
